# lp-ddns/etc/cron.d/lp-ddns
*/30 * * * /usr/local/bin/lp-ddns >/dev/null 2>&1
